# systemd/aetheric-agent.service
[Unit]
Description=Aetheric Edge Agent - MQTT-based edge computing agent
Documentation=https://github.com/Aetheric-Technology/Aetheric-Edge
After=network.target mosquitto.service
Wants=network.target
Requires=mosquitto.service

[Service]
Type=simple
User=aetheric
RuntimeDirectory=aetheric-agent
ExecStartPre=+-/usr/local/bin/aetheric init
ExecStart=/usr/local/bin/aetheric-agent --config /etc/aetheric-edge/aetheric.toml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target